# openrc/init.d/schedule
#!/sbin/openrc-run
# (C) Martin V\"ath <martin@mvath.de>

description='server to schedule and execute queued jobs'

start() {
	ebegin 'Starting schedule-server'
	eval start-stop-daemon -S -b -N 9 \
		-u \"\${schedule_user:-schedule}\" \
		-g \"\${schedule_group:-schedule}\" \
		-- schedule-server ${schedule_opts}
	eend ${?}
}

stop() {
	ebegin 'Stopping schedule-server'
	eval schedule ${schedule_opts} stop-server \
		|| start-stop-daemon -K schedule-server
	eend ${?}
}
