# systemd/system/schedule.service
[Unit]
Description=Server to schedule and execute queued jobs

[Service]
User=root
Group=nobody
SupplementaryGroups=
Nice=19
ExecStart=/usr/bin/schedule-server -Y /etc/schedule.password
ExecStop=/usr/bin/schedule -Y /etc/schedule.password stop-server

[Install]
WantedBy=multi-user.target
