# systemd/system/schedule.socket
[Unit]
Description=schedule-server socket, not reliably usable

[Socket]
ListenStream=127.0.0.1:8471

[Install]
WantedBy=sockets.target
